# systemd/voicebot.service
[Unit]
Description=UIT Prayagraj Voice Chatbot - Alexa-like Voice Assistant
After=network-online.target sound.target pulseaudio.service redis-server.service
Wants=network-online.target pulseaudio.service redis-server.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/chatbot
ExecStart=/home/pi/chatbot/.venv/bin/python -m src.main --no-wake-word
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# ── Audio access ──────────────────────────────────────────────────────
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
Environment=XDG_RUNTIME_DIR=/run/user/1000

# ── General ───────────────────────────────────────────────────────────
Environment=PYTHONUNBUFFERED=1
Environment=WAKE_WORD=ok uit
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
